# scripts/paper-bot.service
[Unit]
Description=Lateralus Paper Bot Queue Runner
After=network.target
# Restart automatically if it crashes

[Service]
Type=simple
User=root
WorkingDirectory=/root/paper_bot
EnvironmentFile=/etc/paper_bot.env
ExecStart=/usr/bin/python3 /root/paper_bot/paper_bot_runner.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/root/paper_bot/logs/service.log
StandardError=append:/root/paper_bot/logs/service_err.log

# Resource limits — only relevant on laptop; on r420 these are generous
CPUQuota=400%
MemoryMax=512M
Nice=10

[Install]
WantedBy=multi-user.target
